# config/etc/supervisor/conf.d/monitor.conf
[program:monitor]
directory=/var/portal
command=python scripts/monitor.py
numprocs=1
autostart=true
autorestart=true
user=pi
