# config/etc/supervisor/conf.d/websockets.conf
[program:websockets]
directory=/var/portal
command=php artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=www-data
