# config/extra/pulser.conf
[program:pulser]
directory=/var/www/portal
command=python scripts/pulse.py
numprocs=1
autostart=true
autorestart=true
user=www-data
