# config/extra/websockets.conf
[program:websockets]
directory=/var/www/portal
command=php artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=www-data
